# deploy/neogogy.service
[Unit]
Description=Neogogy Readiness Index (Next.js)
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/neogogy/app
EnvironmentFile=/opt/neogogy/app/.env.production
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000
Restart=always
RestartSec=5
StandardOutput=append:/var/log/neogogy.log
StandardError=append:/var/log/neogogy.log

[Install]
WantedBy=multi-user.target
